# demo/blawkops-demo.service
[Unit]
Description=BlawkOps Demo API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/BBJJ
ExecStart=/usr/bin/python3 -m uvicorn demo.api:app --host 127.0.0.1 --port 8100
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
